# systemd/amount.path
[Unit]
Description=Wait for directory '/dev/block' to be modified

[Path]
PathModified=/dev/block
Unit=amount.service

[Install]
WantedBy=multi-user.target
