# systemd/automnt.path
[Unit]
Description=Automatically mount removable devices by waiting for new inputs under /dev/block.

[Path]
PathModified=/dev/block
Unit=automnt.service

[Install]
WantedBy=basic.target
