# systemd/automnt.service
[Unit]
Description=Automatically mount inserted devices.

[Service]
Type=oneshot
ExecStart=
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
