# systemd/automnt.service
[Unit]
Description=Automatically mount removable devices.

[Service]
Type=oneshot
ExecStart=

[Install]
WantedBy=basic.target
